# etc/containers/systemd/abcde@.container
[Unit]
Description=Container for processing audio CDs using abcde

[Container]
ContainerName=%p
Image=ghcr.io/cubt85iz/abcde-container:latest
AutoUpdate=registry

[Service]
ExecStopPost=eject %I
