# etc/containers/systemd/ddns-updater.container
[Unit]
Description=Service for updating dynamic DNS address for Cloudflare

[Container]
ContainerName=%p
Image=ddns-updater.build
