# etc/containers/systemd/ddns-updater.container
[Unit]
Description=Container service for updating Cloudflare DNS records

[Container]
ContainerName=%p
Image=ghcr.io/cubt85iz/ddns-updater:latest
