# etc/containers/systemd/embyserver.container
[Unit]
Description=Container service for Emby
Requires=network-online.target nss-lookup.target
After=network-online.target nss-lookup.target

[Container]
ContainerName=%p
Image=docker.io/emby/%p:latest
AddDevice=/dev/dri/renderD128
Volume=${CONTAINER_PATH}/data:/config:Z
PublishPort=${WEB_PORT}:8096
Tmpfs=/transcode
AutoUpdate=registry

[Service]
EnvironmentFile=/etc/environment
Restart=on-failure

[Install]
WantedBy=default.target
