# etc/containers/systemd/immich-redis.container
[Unit]
Description=Container service for Immich Redis Cache
Requires=immich-network.service
After=immich-network.service

[Container]
ContainerName=%p
Image=docker.io/library/redis:7
Network=immich-network
AutoUpdate=registry

[Service]
Restart=on-failure
