# etc/containers/systemd/nextcloud-clamav.container
[Unit]
Description=Container service for Nextcloud AntiVirus
Requires=nextcloud-network.service
After=nextcloud-network.service

[Container]
ContainerName=%p
Image=docker.io/clamav/clamav:latest
Volume=${CONTAINER_PATH}/av:/var/lib/clamav:Z
Network=nextcloud-network

[Service]
ExecCondition=/usr/bin/test -d "${CONTAINER_PATH}/av"
ExecStartPre=/usr/bin/chown -R 100 ${CONTAINER_PATH}/av
Restart=on-failure
