# etc/containers/systemd/stable-diffusion-webui-download.container
[Unit]
Description=

[Container]
Image=ghcr.io

[Service]
Restart=on-failure
RestartSec=30
