# etc/containers/systemd/stable-diffusion-webui-download.container
[Unit]
Description=

[Container]
Image=ghcr.io
Volume=${CONTAINER_PATH}/data:/data:z,rw

[Service]
Restart=on-failure
RestartSec=30
