# etc/systemd/system/backup-yearly@.timer
[Unit]
Description=Backup container volumes yearly

[Timer]
OnCalendar=*-1-1 4:00:00
AccuracySec=1s
