# etc/systemd/system/backup-daily@.timer
[Unit]
Description=Backup container volumes daily

[Timer]
OnCalendar=Mon..Sat *-*-2..31 4:00:00
AccuracySec=1s
