# etc/systemd/system/healthcheck-updates.timer
[Unit]
Description=Check for important updates

[Timer]
OnUnitActiveSec=60m

[Install]
WantedBy=timers.target