# etc/systemd/system/healthcheck-container@.timer
[Unit]
Description=Check health of container every 15 minutes

[Timer]
OnCalendar=*:0/15
AccuracySec=1s

[Install]
WantedBy=timers.target