# etc/systemd/system/healthcheck-zfs.timer
[Unit]
Description=Check health of ZFS pools every 15 minutes

[Timer]
OnBootSec=1m
OnCalendar=*:0/15

[Install]
WantedBy=timers.target