# etc/systemd/system/sync@.timer
[Unit]
Description=Synchronize source & target every minute

[Timer]
OnUnitInactiveSec=1m
AccuracySec=1s
