# etc/systemd/system/ytsubs.timer
[Unit]
Description=Download YouTube videos after rebooting or weekly.

[Timer]
OnBootSec=15min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target